# osquery/osquery.linux.conf
{
    "options": {
            "logger_snapshot_event_type": "true",
            "schedule_splay_percent": 10,
            "logger_path": "osquery/log/",
            "pidfile": "osquery/osquery.pidfile",
            "database_path": "osquery/osquery.db",
            "enable_syslog": "true",
            "enable_monitor": "true"
    },
    "platform": "linux",
    "schedule": {
        "crontab": {
            "query": "SELECT * FROM crontab;",
            "interval": 1100
        },
        "largest_process": {
             "query": "select pid, name, uid, resident_size from processes order by resident_size desc limit 10;",
              "interval": 60
        },
        "system_profile": {
            "query": "SELECT * FROM osquery_schedule;"
        },
        "system_info": {
            "query": "SELECT hostname, cpu_brand, physical_memory FROM system_info;",
            "interval": 3600
        },
        "hardware_events":{
          "query":"select * from hardware_events",
          "interval":5
        },
        "network_connection": {
            "query": "SELECT DISTINCT process.name, listening.port, listening.address, process.pid FROM processes AS process JOIN listening_ports AS listening ON process.pid = listening.pid;",
            "interval":10,
            "description": "Finding new processes listening on network ports"
        },
        "network_activity": {
            "query": "select s.pid, p.name, local_address, remote_address, family, protocol, local_port, remote_port from process_open_sockets s join processes p on s.pid = p.pid where remote_port not in (80, 443) and family = 2;",
            "interval": 10,
            "description": "Finding suspicious outbound network activity - looks for processes with IP traffic to ports not in (80, 443)"
        },
        "kernel_modules": {
            "query": "SELECT * FROM kernel_modules;",
            "interval": 100,
            "description": "Linux kernel modules both loaded and within the load search path."
        }

    },
    "decorators": {
        "load": [
            "SELECT uuid AS host_uuid FROM system_info;",
            "SELECT user AS username FROM logged_in_users ORDER BY time DESC LIMIT 1;"
        ]
    },
    "packs": {
        "fim": "/usr/share/osquery/packs/fim.conf",
        "osquery-monitoring": "/usr/share/osquery/packs/osquery-monitoring.conf",
        "incident-response": "/usr/share/osquery/packs/incident-response.conf",
        "it-compliance": "/usr/share/osquery/packs/it-compliance.conf",
        "vuln-management": "/usr/share/osquery/packs/vuln-management.conf"
    }
}